# etc/udev/rules.d/90-display-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/display-hotplug.sh"
